# packages/123done/ansible/playbooks/roles/123done/files/123done-web.conf
[program:123done-web]
command=node server.js
directory=/home/app/123done
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
stopwaitsecs=3
stdout_logfile=/var/log/123done-web.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/123done-web-err.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
user=app
environment=PORT=10000
