# packages/123done/ansible/playbooks/roles/common/files/supervisor.conf
description "supervisor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/supervisord --configuration /etc/supervisord.conf
